# server/visits.conf
[program:visits]
environment = PYTHONUNBUFFERED=1
directory=/root
command=python3 visits.py
autostart=true
autorestart=true
stderr_logfile=/root/visits.err
stdout_logfile=/root/visits.log
